# config/launchd/com.wechat.articlecrawler.reanalyze-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.wechat.articlecrawler.reanalyze-api</string>

    <key>WorkingDirectory</key>
    <string>/Users/chenwangqian/trae/WeChat-Article-Crawler</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/chenwangqian/trae/WeChat-Article-Crawler/bin/run_reanalyze_api_launchd.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/chenwangqian/trae/WeChat-Article-Crawler/logs/launchd.reanalyze-api.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chenwangqian/trae/WeChat-Article-Crawler/logs/launchd.reanalyze-api.err.log</string>
  </dict>
</plist>
